# health-dashboard/config/com.health-dashboard.streamlit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.health-dashboard.streamlit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/homeserver/Home-Tools/health-dashboard/.venv/bin/streamlit</string>
        <string>run</string>
        <string>dashboard/app.py</string>
        <string>--server.port</string>
        <string>8501</string>
        <string>--server.address</string>
        <string>0.0.0.0</string>
        <string>--server.headless</string>
        <string>true</string>
        <string>--browser.gatherUsageStats</string>
        <string>false</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/homeserver/Home-Tools/health-dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/homeserver/Library/Logs/health-dashboard/streamlit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/homeserver/Library/Logs/health-dashboard/streamlit.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/homeserver</string>
    </dict>
</dict>
</plist>
