# health-dashboard/config/com.health-dashboard.intervals-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.health-dashboard.intervals-poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/homeserver/Home-Tools/health-dashboard/.venv/bin/python3</string>
        <string>-m</string>
        <string>collectors.intervals_collector</string>
        <string>--days</string>
        <string>14</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/homeserver/Home-Tools/health-dashboard</string>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/homeserver/Library/Logs/health-dashboard/intervals-poll.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/homeserver/Library/Logs/health-dashboard/intervals-poll.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>KEYCHAIN_PATH</key>
        <string>/Users/homeserver/Library/Keychains/login.keychain-db</string>
    </dict>
</dict>
</plist>
