# finance-monitor/com.home-tools.finance-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.home-tools.finance-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/homeserver/Home-Tools/finance-monitor/.venv/bin/python3</string>
    <string>/Users/homeserver/Home-Tools/finance-monitor/main.py</string>
    <string>serve</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/homeserver/Home-Tools/finance-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>KEYCHAIN_PATH</key>
    <string>/Users/homeserver/Library/Keychains/login.keychain-db</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/home-tools-finance-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/home-tools-finance-monitor-error.log</string>
</dict>
</plist>
