# event-aggregator/com.antora.event-aggregator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.antora.event-aggregator</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update these paths to match your environment -->
    <string>/usr/local/bin/python3</string>
    <string>/Users/YOURUSERNAME/path/to/Antora-Analytics/event-aggregator/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOURUSERNAME/path/to/Antora-Analytics/event-aggregator</string>

  <!-- Run every 15 minutes -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/tmp/antora-event-aggregator.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/antora-event-aggregator-error.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
