# event-aggregator/com.home-tools.event-aggregator.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.home-tools.event-aggregator.worker</string>

  <!-- Tier 2.4: long-running worker that drains text_queue + ocr_queue
       sequentially and owns model-swap orchestration. Pair with the
       .fetch plist (cron-style) and the dispatcher (event-driven OCR). -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/homeserver/Home-Tools/event-aggregator/.venv/bin/python3</string>
    <string>/Users/homeserver/Home-Tools/event-aggregator/main.py</string>
    <string>worker</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/homeserver/Home-Tools/event-aggregator</string>

  <!-- Long-running — keep restarting if it exits -->
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <!-- Don't restart faster than every 10s if it's crashlooping -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/homeserver/Library/Logs/home-tools/event-aggregator-worker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/homeserver/Library/Logs/home-tools/event-aggregator-worker.log</string>
</dict>
</plist>
